# master_test.service
[Unit]
Description=Go Master Daemon

[Service]
ExecStart=/usr/local/bin/master_test
Restart=always

[Install]
WantedBy=multi-user.target

