# worker_test.service
[Unit]
Description=Go worker Daemon

[Service]
ExecStart=/usr/local/bin/worker_test
Restart=always

[Install]
WantedBy=multi-user.target

